# scripts/com.nyk.docs-maintenance.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.nyk.docs-maintenance</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>REPLACE_WITH_PROJECT_ROOT/scripts/docs_maintenance.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>REPLACE_WITH_PROJECT_ROOT</string>
	<key>StartInterval</key>
	<integer>43200</integer>
	<key>StandardOutPath</key>
	<string>/tmp/docs-maintenance-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/docs-maintenance-stderr.log</string>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
